# run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=local-fs.target

[Mount]
Where=/
What=/dev/disk/by-uuid/da2a890a-eba1-46b8-9b26-dfb8488c3aeb
Type=ext4
Options=errors=remount-ro
